# roci-matrix.service
[Unit]
Description=Roci Matrix Service (Deno)
After=network-online.target roci-agent.service
Requires=roci-agent.service

[Service]
Type=simple
User=tijs
WorkingDirectory=/home/tijs/roci/roci-matrix
EnvironmentFile=/home/tijs/roci/roci-matrix/.env
EnvironmentFile=-/etc/roci/secrets.conf
ExecStart=/home/tijs/.deno/bin/deno task start

Restart=always
RestartSec=10

# Process cleanup
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=30

StandardOutput=journal
StandardError=journal
SyslogIdentifier=roci-matrix

# Resource limits
MemoryMax=512M
CPUQuota=50%

# Security
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
